# mimic.service
# MIMIC Web Server - Systemd Service File
# ==============================================
#
# This service runs the WEB SERVER ONLY (Gunicorn + Flask).
# It does NOT run the Telegram bot or background worker.
# 
# ARCHITECTURE:
#   ✓ mimic.service          → Web server (THIS SERVICE)
#   ✓ mimic-worker.service   → Background task processor (ARQ)
#   ✓ mimic-bot.service      → Telegram bot polling (SEPARATE)
#
# FEATURES:
#   ✓ Gunicorn with eventlet worker (WebSocket support)
#   ✓ Multiple workers for high concurrency
#   ✓ Always running with automatic restart
#   ✓ Graceful shutdown handling
#   ✓ Memory and resource limits
#   ✓ No bot conflicts (bot runs separately)
#
# INSTALLATION:
#   1. Copy this file to: /etc/systemd/system/mimic.service
#   2. Edit paths if needed (default: /var/www/mimic)
#   3. Run: sudo systemctl daemon-reload
#   4. Run: sudo systemctl enable mimic
#   5. Run: sudo systemctl start mimic
#
# COMMANDS:
#   sudo systemctl start mimic      # Start web server
#   sudo systemctl stop mimic       # Stop web server
#   sudo systemctl restart mimic    # Restart web server
#   sudo systemctl status mimic     # Check status
#   journalctl -u mimic -f          # View logs (follow)
#
# DEPLOYMENT ORDER:
#   1. Start mimic.service (web server) ← THIS SERVICE
#   2. Start mimic-worker.service (task processor)
#   3. Start mimic-bot.service (bot polling)
#

[Unit]
Description=MIMIC Web Server (Gunicorn + Flask)
Documentation=https://mimic.cash
After=network-online.target postgresql.service redis.service
Wants=network-online.target postgresql.service redis.service
# Auto-start after boot and network is ready
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=simple
User=mimic
Group=mimic
WorkingDirectory=/var/www/mimic

# Environment
Environment="FLASK_ENV=production"
Environment="PYTHONUNBUFFERED=1"
Environment="PATH=/var/www/mimic/venv/bin:/usr/local/bin:/usr/bin:/bin"
EnvironmentFile=-/var/www/mimic/.env

# Main web server process - Gunicorn with eventlet for WebSockets
# 
# IMPORTANT: This ONLY runs the web server. The Telegram bot runs separately
# via mimic-bot.service to prevent 409 Conflict errors.
#
# BIND ADDRESS OPTIONS:
#   127.0.0.1:8000 - Localhost only (RECOMMENDED with nginx reverse proxy)
#   0.0.0.0:8000   - All interfaces (use if no reverse proxy, or Docker)
#
# For nginx reverse proxy setup (recommended), use 127.0.0.1
# For direct access or Docker, change to 0.0.0.0
#
# WORKER SCALING:
#   You can now safely scale workers (e.g., --workers 4) without bot conflicts!
#   Each worker serves web traffic only. The bot runs independently.
#
ExecStart=/var/www/mimic/venv/bin/gunicorn \
    --worker-class eventlet \
    --workers 1 \
    --threads 100 \
    --worker-connections 1000 \
    --bind 127.0.0.1:8000 \
    --timeout 120 \
    --keep-alive 65 \
    --max-requests 10000 \
    --max-requests-jitter 1000 \
    --graceful-timeout 30 \
    --backlog 2048 \
    --access-logfile /var/www/mimic/logs/access.log \
    --error-logfile /var/www/mimic/logs/error.log \
    --capture-output \
    app:app

# Graceful stop/reload
ExecReload=/bin/kill -s HUP $MAINPID
ExecStop=/bin/kill -s TERM $MAINPID
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# ============ ALWAYS RUNNING - RESTART POLICY ============
# Always restart on failure, exit, or signal
Restart=always
RestartSec=5

# ============ WATCHDOG (Health Monitoring) ============
# Service must ping systemd every 30 seconds or it will be restarted
# Note: Requires app to call sd_notify("WATCHDOG=1") - see below
# WatchdogSec=30

# ============ RESOURCE LIMITS ============
LimitNOFILE=65535
LimitNPROC=4096
# Memory limit (adjust based on your VPS RAM)
MemoryMax=2G
MemoryHigh=1.5G

# ============ SECURITY HARDENING ============
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
PrivateTmp=true
PrivateDevices=true
ProtectKernelTunables=true
ProtectKernelModules=true
ProtectControlGroups=true
ReadWritePaths=/var/www/mimic /var/www/mimic/logs

# ============ LOGGING ============
StandardOutput=append:/var/www/mimic/logs/stdout.log
StandardError=append:/var/www/mimic/logs/stderr.log
SyslogIdentifier=mimic

[Install]
WantedBy=multi-user.target
